# nina/tools/run_script.service
[Unit]
Description=Run Script

[Service]
ExecStart=/udata/nina_test.sh

[Install]
WantedBy=multi-user.target